# scripts/notificatioin services/reminder.timer
[Unit]
Description=Run reminder every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=reminder.service

[Install]
WantedBy=timers.target

